# template/v2/dirs/etc/profile.d/q-cli-message.sh
#!/bin/bash
# Display Q CLI message when terminal starts
echo -e echo -e "\033[1;32m
##############################################################################
#                                                                            #
#                               TERMINAL                                     #
#                         Powered by Amazon Q CLI                            #
#                                                                            #                                                           #
#                  © 2025, Amazon Web Services, Inc. or its affiliates.      #
#                            All rights reserved.                            #
#                                                                            #
#                                                                            #
##############################################################################
\033[0m"